# icsDroneICC/icsdrone-icc/DEBIAN/postinst
#!/bin/sh
update-rc.d icsdrone-icc.sh defaults
